# rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

mkdir /var/run/zabbix
chmod 777 /var/run/zabbix/
/usr/local/sbin/zabbix_agentd

#sleep 60
/home/cubie/wkdir/streamer_run.sh

exit 0
